# contrib/init/connexd.conf
description "CONNEX Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env CONNEXD_BIN="/usr/bin/connexd"
env CONNEXD_USER="connex"
env CONNEXD_GROUP="connex"
env CONNEXD_PIDDIR="/var/run/connexd"
# upstart can't handle variables constructed with other variables
env CONNEXD_PIDFILE="/var/run/connexd/connexd.pid"
env CONNEXD_CONFIGFILE="/etc/connex/connex.conf"
env CONNEXD_DATADIR="/var/lib/connexd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # connexd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$CONNEXD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run connexd."
        echo "The setting must appear in $CONNEXD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"CONNEX Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$CONNEXD_PIDDIR"
    chmod 0755 "$CONNEXD_PIDDIR"
    chown $CONNEXD_USER:$CONNEXD_GROUP "$CONNEXD_PIDDIR"
    chown $CONNEXD_USER:$CONNEXD_GROUP "$CONNEXD_CONFIGFILE"
    chmod 0660 "$CONNEXD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$CONNEXD_PIDFILE" \
    --chuid $CONNEXD_USER:$CONNEXD_GROUP \
    --exec "$CONNEXD_BIN" \
    -- \
    -pid="$CONNEXD_PIDFILE" \
    -conf="$CONNEXD_CONFIGFILE" \
    -datadir="$CONNEXD_DATADIR" \
    -disablewallet \
    -daemon

